# systemd/postfiat-validator-transport.service.example
[Unit]
Description=PostFiat validator transport
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=postfiat
Group=postfiat
ExecStart=/opt/postfiat/bin/postfiat-node transport-validator-serve --unsafe-devnet-json-storage --unsafe-devnet-file-signer --data-dir /var/lib/postfiat/validator-0 --topology /etc/postfiat/private-overlay-topology.json --key-file /var/lib/postfiat/validator-0/validator-keys.json --vote-dir /var/lib/postfiat/validator-0/votes --bind-host 10.0.0.10 --max-connections 10000 --timeout-ms 30000 --event-log /var/log/postfiat/transport-events.jsonl
RuntimeDirectory=postfiat-validator-transport
RuntimeDirectoryMode=0700
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=control-group
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/var/lib/postfiat/validator-0 /var/log/postfiat /run/postfiat-validator-transport
UMask=0077
LimitNOFILE=65536
LimitCORE=0
TasksMax=1024

[Install]
WantedBy=multi-user.target
